# build/ansible/playbooks/roles/runner/files/bb-runner.service
[Unit]
Description=Runs the BuildBeaver Runner process
After=network.target bb-ephemeral-ssd.service

[Service]
ExecStart=/usr/bin/start-bb-runner

[Install]
WantedBy=multi-user.target